# etc/init.d/weewx
#!/sbin/openrc-run

name="weewx"

: ${command_user:="weewx"}
: ${start_wait=100}  # milliseconds
: ${retry="TERM/10/KILL/5"}

command="/home/weewx/weewx-venv/bin/weewxd"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

start_stop_daemon_args="--wait $start_wait $start_stop_daemon_args"
# The leading space is to avoid fallback to $start_stop_daemon_args when this
# is empty (supervise-daemon doesn't support --wait).
supervise_daemon_args=" $supervise_daemon_args"

export WEEWX_DATA="${datadir:="/home/weewx/weewx-data"}"

depend() {
	need dev localmount
	use mosquitto
}
